# runssserver.conf
start on (runlevel [2345])
stop on (runlevel [016])

pre-start script
/etc/init.d/runssserver start
end script

post-stop script
/etc/init.d/runssserver stop
end script